# deployment/django-admin-automation.service
[Unit]
Description=Django admin automation dashboard (Gunicorn)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/django_admin_automation
Environment=PATH=/opt/django_admin_automation/venv/bin
Environment=FLASK_DEBUG=0
EnvironmentFile=-/etc/django-admin-automation.env
# Prefer run_production.sh so .secrets.env + secrets.local.env + .secrets.key load the same as local.
ExecStart=/opt/django_admin_automation/run_production.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
